# supervisord.conf
[supervisord]
nodaemon=true

[program:back]
command=/home/back/gradlew clean build bootRun
autorestart=true
directory=/home/back
